# mcp-deployment/systemd/1c-mcp-metadata.service
[Unit]
Description=1C MCP Code Metadata Search Server
After=docker.service network-online.target
Wants=docker.service network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/mcp
ExecStart=/usr/local/bin/docker compose up -d 1c-metadata
ExecStop=/usr/local/bin/docker compose stop 1c-metadata

[Install]
WantedBy=multi-user.target
